# conf/cron.d/logrotate_init
0 3 * * 7 root /usr/sbin/logrotate -f /etc/logrotate.conf
