# conf/cron.d/php_init
0 4 * * * root /usr/bin/find /var/lib/php/sessions -mtime +5 -exec rm -f {} \;
30 4 * * * root /usr/bin/find /tmp/php/sessions -mtime +5 -exec rm -f {} \;
